# roles/locust_master/templates/locustmaster.conf.j2
[program:locustmaster]
command=/usr/local/bin/locust -f /opt/app/agent.py --master --host={{ test_host }}
redirect_stderr=true
stdout_logfile=/var/log/locust_master.log
stdout_logfile_maxbytes=100MB
